# deepin-big-improvements/DEBIAN/postinst
#! /bin/sh

set -e

systemctl daemon-reload
systemctl enable bigswappiness.service 


exit 0
